# logs/etc/systemd/system/wg-quick@wg0.service.d/logging.conf
[Service]
ExecStartPost=/bin/bash -c 'journalctl -u wg-quick@%i -f >> /var/log/wireguard/%i.log & echo $! > /run/wg-quick-%i-log.pid'
ExecStopPost=/bin/bash -c 'if [ -f /run/wg-quick-%i-log.pid ]; then PID=$(cat /run/wg-quick-%i-log.pid); if [ -n "$PID" ] && kill -0 $PID 2>/dev/null; then kill $PID && rm /run/wg-quick-%i-log.pid; fi; fi'